# vagrant.sndcomp/sndcomp.upstart
# https://github.com/lingdb/container/issues/1
# https://unix.stackexchange.com/a/84289/23749

description "Run '/containerSetup/vagrant.sndcomp/run.sh'."
author "Jakob Runge <mam09crm@studserv.uni-leipzig.de>"

start on runlevel [2345]
stop on runlevel [016]

chdir /containerSetup/vagrant.sndcomp

exec /bin/bash run.sh
