# database/dumpdb.timer
[Unit]
Description=Run dumpdb.service every day.

[Timer]
OnBootSec=15min
OnUnitActiveSec=1day
Unit=dumpdb.service

[Install]
WantedBy=timers.target
