# install/ielex2.upstart
# Inspired by soundcomparisons.upstart

description "Run '/containerSetup/ielex2/manage.py runserver'."
author "Jakob Runge <mam09crm@studserv.uni-leipzig.de>"

start on runlevel [2345]
stop on runlevel [016]

respawn

chdir /containerSetup/ielex2

exec ./bin/python manage.py runserver 127.0.0.1:5001
