# flask/flask.service
[Unit]
Description=Run '/containerSetup/soundcomparisons/make main'.
Requires=mysql.service

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/make -C /containerSetup/soundcomparisons main
User=vagrant
Group=systemd-journal
Restart=always
RestartSec=30
